# components/example/component/applications/savory-example/crontab
*/1 * * * * <% document.executeOnce('/savory/service/cache/'); new Savory.Cache({name: 'result', logLevel: 'info'}).prune(); %>

